# socketrelayserver.service
[Unit]
Description=SocketRelay Service
After=network.target

[Service]
WorkingDirectory=/etc/socketrelayserver
ExecStart=/etc/socketrelayserver/run.php
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
User=nobody

[Install]
WantedBy=multi-user.target
